# ops/cron/youtube-curator.cron
0 8 * * * cd /home/dan-parii/Documents/HermesYoutubeCurator && ./scripts/morning_run.py
