# ops/systemd/hermes-youtube-curator.service
[Unit]
Description=Hermes YouTube Curator

[Service]
Type=oneshot
WorkingDirectory=/home/dan-parii/Documents/HermesYoutubeCurator
ExecStart=/usr/bin/env python3 scripts/morning_run.py
